# files/sentry/sentry-cron.service
[Unit]
Description=Sentry Beat Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/.venv/sentry
Environment=SENTRY_CONF=/home/ubuntu/sentry
ExecStart=/home/ubuntu/.venv/sentry/bin/sentry run cron

[Install]
WantedBy=multi-user.target
